# harness/soak-ab/soak-b.service
# Soak B — Portal-native realtime soak (mirror of Soak A, Portal /stream instead of RPC realtime).
# RENDERED by harness/soak-ab/deploy-soak-b.sh (the @@…@@ tokens are substituted); do not install by hand.
#
# GUARDRAILS baked in: dedicated DB euler_rt_b + dedicated port (never the euler prod DB or :9547),
# secrets only via a chmod-600 EnvironmentFile, RSS-capped, auto-restart. The unit is enabled but
# NOT started by the deploy script — the operator starts it deliberately.

[Unit]
Description=portal-ponder Soak B (Portal-native realtime, DB euler_rt_b)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple

# Run as an UNPRIVILEGED account, never root. Rendered by deploy-soak-b.sh from $SOAK_B_USER
# (default: the invoking operator's own login, resolved at deploy time — never hard-coded here).
User=@@USER@@
Group=@@GROUP@@

WorkingDirectory=@@WORKDIR@@

# All secrets + tunables live in this chmod-600 file (PORTAL_API_KEY, SQD_RPC_KEY, DATABASE_URL, …).
EnvironmentFile=@@ENVFILE@@

# Portal-native realtime + on-by-default runtime invariant checks (a violation crashes loud).
Environment=PORTAL_REALTIME=stream
Environment=PORTAL_CHECKS=on
Environment=PONDER_LOG_LEVEL=info
Environment=CI=true

# The soak indexes the SAME chains as Soak A. The app's ponder.config.ts reads EULER_CHAINS (full
# chains.json names, e.g. ethereum,base,arbitrum) — written into the EnvironmentFile above by
# deploy-soak-b.sh from the SOAK_CHAINS/EULER_CHAINS knob (validated against chains.json). It is NOT
# hard-coded here so a redeploy against a differently-scoped soak stays consistent from the env file.

# Restart counting is a first-class soak signal. Each (re)start appends a UTC timestamp to the
# restart log; ab-diff.mjs reads it into restartCount / lastRestartAt and alerts on >3 restarts/hour
# (crash-loop). In PORTAL_REALTIME=stream mode unknown-head and reorg-gap conditions are FATAL by
# design and exit 75 (EX_TEMPFAIL) — Restart=on-failure makes that the intended restart-recovery,
# NOT a soak failure (only the crash-loop RATE is an alert).
Environment=SOAK_B_RESTART_LOG=@@RESTART_LOG@@
ExecStartPre=/bin/sh -c 'echo "$(date -u +%%Y-%%m-%%dT%%H:%%M:%%SZ) restart" >> "$SOAK_B_RESTART_LOG"'

# --schema is rendered from the single SOAK_B_SCHEMA knob (default soak_b) so it stays consistent
# with DATABASE_SCHEMA in the EnvironmentFile — a redeploy over an existing soak using a different
# schema must not silently point the app at soak_b. Never hard-code the schema here.
ExecStart=@@WORKDIR@@/node_modules/.bin/ponder start --schema @@SCHEMA@@ --port @@PORT@@

# RSS cap — MemoryHigh throttles, MemoryMax is the hard kill ceiling (the soak asserts it holds).
MemoryHigh=@@MEM_HIGH@@
MemoryMax=@@MEM_MAX@@

# on-failure covers exit 75 (the designed stream-mode restart) and any crash/OOM/signal.
Restart=on-failure
RestartSec=10
# 2 planned SIGKILL+restart drills mid-soak are exercised by the operator via `systemctl kill -s KILL`.
KillSignal=SIGTERM
TimeoutStopSec=30

# Hardening: unprivileged, no new privileges, minimal filesystem + kernel surface. The process only
# needs its workdir and the restart-log dir writable; everything else is read-only or hidden.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=@@WORKDIR@@ @@RESTART_LOG_DIR@@
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
